# conf/supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)
chown=ubuntu

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf

[program:cadre-start-jg-listener]
environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:%(ENV_PATH)s
user=ubuntu
command = /home/ubuntu/cloud-admin/deployment/start-cluster.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/cadre-start-cluster.super.err.log
stdout_logfile=/var/log/cadre-start-cluster.super.out.log

[program:janusgraph]
user=ubuntu
directory=/home/ubuntu/janusgraph
redirect_stderr=true
stdout_logfile=/home/ubuntu/janusgraph/logs/gremlin-server.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241,JANUSGRAPH_HOME=/home/ubuntu/janusgraph
stopasgroup=true
command=/home/ubuntu/janusgraph/bin/gremlin-server.sh conf/gremlin-server/socket-http-gremlin-server.yaml

[program:listener1]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 1 -Xms512M -Xmx3840M

[program:listener2]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 2 -Xms512M -Xmx3840M

[program:listener3]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 3 -Xms512M -Xmx3840M

[program:listener4]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 4 -Xms512M -Xmx3840M

[program:listener5]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 5 -Xms512M -Xmx3840M

[program:listener6]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 6 -Xms512M -Xmx3840M

[program:listener7]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 7 -Xms512M -Xmx3840M

[program:listener8]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 8 -Xms512M -Xmx3840M

[program:listener9]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 9 -Xms512M -Xmx3840M

[program:listener10]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 10 -Xms512M -Xmx3840M

[program:listener11]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 11 -Xms512M -Xmx3840M

[program:listener12]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 12 -Xms512M -Xmx3840M

[program:listener13]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 13 -Xms512M -Xmx3840M

[program:listener14]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 14 -Xms512M -Xmx3840M

[program:listener15]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 15 -Xms512M -Xmx3840M

[program:listener16]
user=ubuntu
directory=/home/ubuntu/janus-job-listener
redirect_stderr=true
stdout_logfile=/home/ubuntu/logs/janus-listener.log
environment=JAVA_HOME=/home/ubuntu/software/jdk1.8.0_241
command=java -cp /home/ubuntu/lib/janus-job-listener-0.0.1-SNAPSHOT-jar-with-dependencies.jar iu.cadre.listeners.job.JobListener /home/ubuntu/lib/cadre_config.properties 16 -Xms512M -Xmx3840M

